# update.service
[Unit]
Description=My script

[Service]
ExecStart=/home/dns/run.sh