# update.service
[Unit]
Description=Godaddy DDNS IP update

[Service]
Type=oneshot
ExecStart=/media/dns/update
